# deploy/kilnwatch.service
[Unit]
Description=KilnWatch thermocouple monitor and web server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/kilnwatch
ExecStart=/home/pi/kilnwatch/venv/bin/python -m kilnwatch.main
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
